# workers/supervisor_tracking_worker.conf
[program:viax_tracking_worker]
command=php /var/www/viax/backend/workers/tracking_stream_worker.php
autostart=true
autorestart=true
numprocs=2
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/log/viax_tracking_worker.log
stderr_logfile=/var/log/viax_tracking_worker_error.log
stopasgroup=true
killasgroup=true
